# scripts/com.kemmer.df-kpm-portfolio-tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- DF-KPM-Portfolio-Tracker LaunchAgent [CRUX-MK] -->
<!-- K_0-MAX: RunAtLoad=false (Phronesis-Pflicht vor Aktivierung) -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kemmer.df-kpm-portfolio-tracker</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/make/Projects/dark-factories/df-kpm-portfolio-tracker/scripts/run-df-kpm-portfolio-tracker-mac.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- K_0-Schutz: Nicht beim Login starten, Phronesis-Pflicht -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/make/.df-kpm-portfolio-tracker/launchagent.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/make/.df-kpm-portfolio-tracker/launchagent.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>DF_KPM_PORTFOLIO_REAL_ENABLED</key>
    <string>false</string>
  </dict>
</dict>
</plist>
